# systemd/farmacafe-monitor.timer.template
# Install/update this timer into ~/.config/systemd/user/ with:
#   ./bin/install_systemd_monitor.sh
[Unit]
Description=Run Cafeteria Farmacia UV watcher on weekdays every 10 minutes (08:01-14:51)

[Timer]
Unit=farmacafe-monitor.service
OnCalendar=Mon..Fri *-*-* 8,9,10,11,12,13,14:01/10:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
